# knepp/cloud-backup.service
[Unit]
Description=Use with timer to back up audio files with rclone

[Service]
ExecStart=/home/data/alicee/knepp_extra/cloud-backup.sh
